# bundle/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:dashboard-backend]
command=node /opt/dashboard-backend/server.js
autostart=true
autorestart=true
environment=PORT="3005"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:calendar-backend]
command=/usr/local/bin/run-backend.sh /usr/local/bin/calendar-backend
autostart=true
autorestart=true
environment=PORT="3006",DATA_PATH="/data/calendar",RUST_LOG="info",HOME="/root",ADMIN_PASSWORD="%(ENV_ADMIN_PASSWORD)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:expense-backend]
command=/usr/local/bin/run-backend.sh /usr/local/bin/expense-tracker-backend
autostart=true
autorestart=true
environment=PORT="3002",DATA_PATH="/data/expenses",RUST_LOG="info",HOME="/root",ADMIN_PASSWORD="%(ENV_ADMIN_PASSWORD)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:inventory-backend]
command=/usr/local/bin/run-backend.sh /usr/local/bin/inventory-backend
autostart=true
autorestart=true
environment=PORT="3003",DATA_PATH="/data/inventory",RUST_LOG="info",HOME="/root",ADMIN_PASSWORD="%(ENV_ADMIN_PASSWORD)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:habits-backend]
command=/usr/local/bin/run-backend.sh /usr/local/bin/habits-backend
autostart=true
autorestart=true
environment=PORT="3010",DATA_PATH="/data/habits",RUST_LOG="info",HOME="/root",ADMIN_PASSWORD="%(ENV_ADMIN_PASSWORD)s",TASKBOARD_URL="http://127.0.0.1:3009",HABITS_URL="http://127.0.0.1:3010",OCTIRON_API_KEY_PATH="%(ENV_OCTIRON_API_KEY_PATH)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:mood-backend]
command=/usr/local/bin/run-backend.sh node /opt/mood-backend/server.js
autostart=true
autorestart=true
environment=PORT="3004",DATA_PATH="/data/mood",NODE_ENV="production",ADMIN_PASSWORD="%(ENV_ADMIN_PASSWORD)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:task-backend]
command=/usr/local/bin/run-backend.sh node /opt/task-backend/server.js
autostart=true
autorestart=true
environment=PORT="3009",DATA_PATH="/data/tasks",NODE_ENV="production",ADMIN_PASSWORD="%(ENV_ADMIN_PASSWORD)s",HABITS_URL="http://127.0.0.1:3010",HABITS_SYNC_USER="aldus",OCTIRON_API_KEY_PATH="%(ENV_OCTIRON_API_KEY_PATH)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
